# packaging/kylin/wps-tts.service
[Unit]
Description=WPS offline read-aloud service
After=sound.target

[Service]
Type=simple
ExecStart=/opt/wps-read-aloud/daemon/wps-tts-daemon -config /etc/wps-read-aloud/config.yaml
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
